# robot_pi/systemd/module-robot-bringup.service
[Unit]
Description=Module Robot manual-safe ROS 2 bringup
Documentation=file://@ROBOT_ROOT@/robot_pi/README.md
# Keep safety/gateway alive in DISCONNECTED while the bridge service restarts.
Wants=module-robot-bridge.service
After=network-online.target module-robot-bridge.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=@ROBOT_USER@
WorkingDirectory=@ROBOT_ROOT@/robot_pi/ros2_ws
EnvironmentFile=/etc/module-robot/runtime.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && source @ROBOT_ROOT@/robot_pi/ros2_ws/install/setup.bash && exec ros2 launch module_robot_bringup manual_bringup.launch.py start_bridge:=false start_gateway:=true autonomous:=false safety_config:=@ROBOT_ROOT@/robot_pi/config/safety.yaml gateway_config:=@ROBOT_ROOT@/robot_pi/config/network.yaml'
Restart=on-failure
RestartSec=2
KillSignal=SIGINT
KillMode=control-group
TimeoutStopSec=8
SendSIGKILL=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=module-robot-bringup
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
